# conf/magento.cron
*/5 * * * * sh /home/magento/files/html/cron.sh >> /home/magento/files/html/var/log/magento.cron.log
